# flask-radius/s6/etc/services.d/update-user/run
#!/usr/bin/with-contenv sh

while :
do
    cd /root/
    sleep 10
    flask update-userdict
    sleep 7200
done
